# initial-conf/crontab
* * * * * root /sbin/pfctl -t blocked_sites -T expire 60
